# nginx-admin-docker/nginx-admin-agent/qas/debian/build/supervisord.conf
[supervisord]
nodaemon=true
	
[program:nginx-agent]
command = /opt/nginx-agent-2.0.1/scripts/debian/nginx-agent.sh start