# nginx-admin-docker/nginx-admin-ui/qas/red-hat/build/supervisord.conf
[supervisord]
nodaemon=true
	
[program:nginx-admin]
command = /opt/nginx-admin/scripts/red-hat/nginx-admin.sh start