# nginx-admin-docker/qas/nginx-admin/debian/build/supervisord.conf
[supervisord]
nodaemon=true
	
[program:nginx-admin]
command = /opt/nginx-admin-2.0.2/scripts/debian/nginx-admin.sh start