# nginx-admin-docker/nginx-admin-agent/qas/red-hat/build/supervisord.conf
[supervisord]
nodaemon=true
	
[program:nginx-agent]
command = /opt/nginx-agent/scripts/red-hat/nginx-agent.sh start