# nginx-admin-docker/qas/nginx-admin-agent/debian/build/supervisord.conf
[supervisord]
nodaemon=true
	
[program:nginx-agent]
command = /opt/nginx-admin-agent-2.0.2/scripts/debian/nginx-admin-agent.sh start