# nginx-admin-docker/qas/nginx-admin-agent/red-hat/build/supervisord.conf
[supervisord]
nodaemon=true
	
[program:nginx-agent]
command = /opt/nginx-admin-agent-2.0.1/scripts/red-hat/nginx-admin-agent.sh start