# nginx-admin-docker/qas/red-hat/build/supervisord.conf
[supervisord]
nodaemon=true
		
[program:nginx-agent]
command=runuser nginx-agent -c "/usr/bin/java -jar -Djava.net.preferIPv4Stack=true /opt/nginx-agent/bin/nginx-admin-agent-standalone-1.0.7-swarm.jar -c /opt/nginx-agent/conf/nginx-agent.conf" >> /opt/nginx-agent/console.log 2>&1 & echo $! > /var/run/nginx-agent/nginx-agent.pid
		
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"


