# nginx-admin-docker/qas/red-hat/build/supervisord.conf
[supervisord]
nodaemon=true
	
[program:nginx-agent]
command=/usr/bin/java -jar -Djava.net.preferIPv4Stack=true /opt/nginx-agent/bin/nginx-admin-agent-standalone-1.0.7-swarm.jar -c /opt/nginx-agent/conf/nginx-agent.conf
		
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"


